# server/pwgen.service
[Unit]
Description=pwgen-server

[Service]
ExecStart=/var/www/pwgen-react/server/start.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target

